# teamviewer/root/etc/services.d/teamviewerd/run
#!/usr/bin/with-contenv bash

/opt/teamviewer/tv_bin/teamviewerd -d
