# teamviewer/root/etc/services.d/teamviewerd/run
#!/usr/bin/execlineb -P

/opt/teamviewer/tv_bin/teamviewerd -f
