# nzbget/root/etc/services.d/nzbget/run
#!/usr/bin/execlineb -P
# vim:set ft=sh sw=2 sts=2 st=2 et:

with-contenv
multisubstitute
{
  importas -D   1000       APP_UID         APP_UID
  importas -D   1000       APP_GID         APP_GID
  importas -D   "appuser"  APP_USER        APP_USER
  importas -D   "/config"  APP_GCONFIG     APP_GCONFIG
  importas -D   "002"      UMASK           UMASK
  importas                 APP_NAME        APP_NAME
}

define HOME "/home/${APP_USER}"
define APP_DIR "/usr/lib/${APP_NAME}"

cd ${APP_DIR}
s6-applyuidgid -u ${APP_UID} -g ${APP_GID} umask ${UMASK} ${APP_DIR}/${APP_NAME} -s -o OutputMode=log -c ${APP_GCONFIG}/${APP_NAME}.conf
