# sonarr/root/etc/services.d/sonarr/run
#!/usr/bin/execlineb -P
# vim:set ft=sh sw=2 sts=2 st=2 et:

with-contenv
multisubstitute
{
  importas -D   1000        APP_UID        APP_UID
  importas -D   1000        APP_GID        APP_GID
  importas -D   "appuser"   APP_USER       APP_USER
  importas -D   "/config"   APP_GCONFIG    APP_GCONFIG
  importas -D   "002"       UMASK          UMASK
  importas                  APP_NAME       APP_NAME
  define MONO_THREADS_PER_CPU 250
  define MONO_GC_PARAMS nursery-size=128m
}

define HOME "/home/${APP_USER}"
define APP_DIR "/opt/${APP_NAME}"
define SONARR "/opt/${APP_NAME}/NzbDrone.exe"

cd $APP_DIR
s6-applyuidgid -u ${APP_UID} -g ${APP_GID} umask ${UMASK} /usr/bin/mono-sgen ${SONARR} \
  -nobrowser \
  -data=${APP_GCONFIG}
